# docker/etc/cron.d/hsn2-thug-docker-disk-guardian
*/15 * * * * root /usr/bin/hsn2-thug-docker-restarter 0.8
